# src/root/Library/Lockdown/ServiceAgents/com.apple.afc2.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
	<dict>
		<key>Label</key>
		<string>com.apple.afc2</string>
		<key>ProgramArguments</key>
		<array>
			<string>/usr/libexec/afcd</string>
			<string>--lockdown</string>
			<string>-S</string>
			<string>-d</string>
			<string>/</string>
		</array>
	</dict>
</plist>
